# server-script/aihomeworkstaging.service
[Unit]
Description=AIHomework-Dev
After=network.target

[Service]
WorkingDirectory=/var/www/aihomework/boilerplate/staging
ExecStart=/bin/bash /var/www/aihomework/boilerplate/staging/server-script/startappstaging.sh
#Restart=on-failure
#RestartSec=20s
StartLimitInterval=0

[Install]
WantedBy=multi-user.target
